# master/deployment/ffmpeg-master.service
[Unit]
Description=FFmpeg RTMP Master Coordination Service
Documentation=https://github.com/psantana5/ffmpeg-rtmp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ffmpeg
Group=ffmpeg
WorkingDirectory=/opt/ffmpeg-rtmp

# Environment: Set API key (REQUIRED for production)
Environment="MASTER_API_KEY=your-secure-api-key-here-change-me"

# Start master service with production defaults
# TLS enabled by default, SQLite persistence, retry logic
ExecStart=/opt/ffmpeg-rtmp/bin/master \
    --port 8080 \
    --db /var/lib/ffmpeg-rtmp/master.db \
    --tls \
    --cert /etc/ffmpeg-rtmp/certs/master.crt \
    --key /etc/ffmpeg-rtmp/certs/master.key \
    --max-retries 3 \
    --metrics \
    --metrics-port 9090

# Restart policy
Restart=always
RestartSec=10s
StartLimitBurst=5
StartLimitIntervalSec=300

# Resource limits
MemoryLimit=2G
MemoryAccounting=yes
CPUQuota=200%
CPUAccounting=yes

# Security hardening
NoNewPrivileges=true
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/ffmpeg-rtmp /var/lib/ffmpeg-rtmp /etc/ffmpeg-rtmp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ffmpeg-master

# Graceful shutdown
TimeoutStopSec=30s
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
